# systemd/peetsfea-web.service
[Unit]
Description=Peets FEA Runner Web Status API (user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Projects/PythonProjects/peetsfea-runner
Environment=PEETSFEA_DB_PATH=%h/Projects/PythonProjects/peetsfea-runner/peetsfea_runner.duckdb
Environment=PEETSFEA_WEB_HOST=127.0.0.1
Environment=PEETSFEA_WEB_PORT=8765
ExecStart=%h/Projects/PythonProjects/peetsfea-runner/.venv/bin/python -c "from peetsfea_runner.systemd_web import run_status_server; run_status_server()"
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
